# wardrobe/com.wardrobe.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wardrobe.app</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/你的用户名/wardrobe && source .venv/bin/activate && exec python app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/你的用户名/wardrobe</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/你的用户名/wardrobe/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/你的用户名</string>
        <key>WARDROBE_MODELS_DIR</key>
        <string>/Users/你的用户名/wardrobe_models</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/你的用户名/wardrobe/logs/wardrobe.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/你的用户名/wardrobe/logs/wardrobe.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
